# src/frontend/startkde/plasma-login-wayland.target
[Unit]
Requires=plasma-login-kwin_wayland.service
BindsTo=plasma-login-kwin_wayland.service
Requires=plasma-login.service
BindsTo=plasma-login.service

Wants=plasma-wallpaper.service
Wants=plasma-kded6.service
